# transportMinistry_backend/interface_api.conf
[program:interface_backend]
directory=/home/bruno1993/transport_api
command= python /home/bruno1993/transport_api/interface_api.py &
autostart=true
autorestart=true